# etc/init.d/nspchSecurityBot
#!/sbin/openrc-run
# Copyright 2022 I'U.R.S.S People's Court for Human Rights (People's Tribunal).
# Distributed under the terms of the GNU General Public License v2

start() {
	ebegin "Starting ${TBOT_NAME}"
	mkdir -p "${TBOT_PIDS}" || return 1
	touch $TBOT_PIDS/$TBOT_NAME.timestamp
	chown $TBOT_USER:$TBOT_GROUP $TBOT_PIDS/$TBOT_NAME.timestamp
	source "$TBOT_VENV"
	start-stop-daemon --start --background --quiet \
	                  --user "${TBOT_USER}" \
	                  --group "${TBOT_GROUP}" \
	                  --make-pidfile --pidfile "${TBOT_PIDS}/${TBOT_NAME}.pid" \
	                  --chdir "${TBOT_DIR}/${TBOT_NAME}" \
	                  --name "${TBOT_NAME}" \
	                  ${TBOT_DIR}/${TBOT_NAME}/${TBOT_NAME}.py \
	                  -- ${TBOT_ARGS}
	eend ${?}
}

stop() {
	ebegin "Stopping ${TBOT_NAME}"
	start-stop-daemon --stop --quiet --retry=TERM/30/KILL/5 --pidfile "${TBOT_PIDS}/${TBOT_NAME}.pid"
	#pkill -P `cat ${PIDFILE}`
	#RETVAL="$?"
	RETVAL=0
	#rm -f ${PIDFILE}
	eend ${RETVAL}
}

